# minecraftjobs.cron
*/15 * * * * ftbuser /bin/bash service minecraft backup
